# system-scripts/systemd-units/zpool-scrub@.service
[Unit]
Description=ZFS scrub on %i

[Service]
Type=oneshot
ExecStart=/usr/bin/zpool scrub %i
